# supervisord.conf
[supervisord]
nodaemon=true

[program:tiarra]
command=/tiarra/tiarra -d --config=/tiarra/tiarra.conf -d
