# supervisord.conf
[supervisord]
nodaemon=true

[program:tiarra]
command=/tiarra/tiarra -d --config=/tiarra/tiarra.conf -d

[program:ruby]
command=ruby /atig/exe/atig -d
